# tools/portfolio/master-kanban/master-kanban-serve-staging.service
[Unit]
Description=Master-Kanban serve STAGING (mario-brain Portfolio API + static cockpit)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/stack/bin/master-kanban serve --port 7770
EnvironmentFile=/etc/master-kanban/db.staging.env
Environment=PORTFOLIO_API_KEY=dev-secret
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target
